# dovecot.service.in
[Unit]
Description=Dovecot IMAP/POP3 email server
After=local-fs.target network.target

[Service]
Type=simple
ExecStart=@sbindir@/dovecot -F
NonBlocking=yes

[Install]
WantedBy=multi-user.target
